# debian/cron.daily
#!/bin/sh
#
# netcomics daily cron job.
#

# Add here calls to netcomics:
#netcomics -c uf -n 1 -g
